# systemd/papergate.service
[Unit]
Description=paperGate e-Paper Display Manager
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/paperGate/core
ExecStart=/usr/bin/python3 /home/pi/paperGate/core/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
